# deploy/aws/supervisor-queue-worker.conf
; Install: sudo cp deploy/aws/supervisor-queue-worker.conf /etc/supervisor/conf.d/flyertrade-queue.conf
; Then: sudo supervisorctl reread && sudo supervisorctl update && sudo supervisorctl start flyertrade-queue:*
;
; Replace /var/www/flyertrade with your project path and www-data with your PHP-FPM user.

[program:flyertrade-queue]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/flyertrade/artisan queue:work database --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/flyertrade/storage/logs/queue-worker.log
stopwaitsecs=3600
